# tls-server/tls-fingerprint.service
    [Unit]
Description=TLS Fingerprint Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tls-fingerprint
ExecStart=/opt/tls-fingerprint/tls-server -port 8443 -cert /opt/tls-fingerprint/server.crt -key /opt/tls-fingerprint/server.key
Restart=always
RestartSec=5

# 安全设置
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tls-fingerprint

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tls-fingerprint

[Install]
WantedBy=multi-user.target
